# types/laravel/docker/conf/cron/crontab
* * * * * cd /var/www/html && php artisan schedule:run 2>&1 >> /var/www/html/storage/logs/schedule.log
